# aitekeji/守护进程配置/atms-task.ini
[program:atms-task]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/atms-task/atms-task.jar --spa=test --spring.profiles.active=test
directory   = /usr/local/atms/atms-task/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/atms-task.log
stderr_logfile=/var/log/atms-task.err.log
