# aitekeji/守护进程配置/atms-trade.ini
[program:atms-trade]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/atms-trade/atms-trade.jar --spa=simulate --spring.profiles.active=simulate
directory   = /usr/local/atms/atms-trade/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/atms-trade.log
stderr_logfile=/var/log/atms-trade.err.log
