# aitekeji/守护进程配置/atms-wallet.ini
[program:atms-wallet]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/atms-wallet/atms-wallet.jar --spa=test --spring.profiles.active=test
directory   = /usr/local/atms/atms-wallet/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/atms-wallet.log
stderr_logfile=/var/log/atms-wallet.err.log
