# aitekeji/守护进程配置/club-admin.ini
[program:club-admin]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/club-admin/club-admin.jar --spa=simulate --spring.profiles.active=simulate
directory   = /usr/local/atms/club-admin/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/club-admin.log
stderr_logfile=/var/log/club-admin.err.log
