# aitekeji/守护进程配置/club-web.ini
[program:club-web]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/club-web/club-web.jar --spa=test --spring.profiles.active=test
directory   = /usr/local/atms/club-web/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/club-web.log
stderr_logfile=/var/log/club-web.err.log
