# aitekeji/守护进程配置/eth-club-backend.ini
[program:eth-club-backend]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/eth-club-backend/eth-club-backend.jar --spa=simulate --spring.profiles.active=simulate
directory   = /usr/local/atms/eth-club-backend/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/eth-club-backend.log
stderr_logfile=/var/log/eth-club-backend.err.log
