# aitekeji/守护进程配置/timing-task.ini
[program:timing-task]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/timing-task/timing-task.jar --spa=simulate --spring.profiles.active=simulate
directory   = /usr/local/atms/timing-task/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/timing-task.log
stderr_logfile=/var/log/timing-task.err.log
