# aitekeji/守护进程配置/tx-manager.ini
[program:tx-manager]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/tx-manager/tx-manager.jar --spa=simulate --spring.profiles.active=simulate
directory   = /usr/local/atms/tx-manager/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/tx-manager.log
stderr_logfile=/var/log/tx-manager.err.log
