# aitekeji/守护进程配置/test/atms-issue.ini
[program:atms-issue]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/atms-issue/atms-issue.jar --spa=test --spring.profiles.active=test
directory   = /usr/local/atms/atms-issue/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/vdb/superlog/atms-issue.log
stderr_logfile=/vdb/superlog/atms-issue.err.log
