# aitekeji/守护进程配置/test/eth-club.ini
[program:eth-club]

command     = java -jar -Xms256m -Xmx512m /usr/local/atms/eth-club/eth-club.jar --spa=test --spring.profiles.active=test
directory   = /usr/local/atms/eth-club/
user        = root
startsecs   = 3


autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/vdb/superlog/eth-club.log
stderr_logfile=/vdb/superlog/eth-club.err.log
